# crontab/haya
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp//crontab.CoJus18781 installed on Wed Jul 20 13:35:51 2005)
# (Cron version V5.0 -- $OpenBSD: crontab.c,v 1.48 2005/01/30 21:00:31 millert Exp $)
#minute hour    mday    month   wday    command
*/10	*	*	*	*	/home/will/bin/getmail
15	*/4	*	*	*	/home/will/Backups/mail.sh	
