# crontab/hayabusa
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp//crontab.SQsoegD855 installed on Wed Jul 27 17:35:44 2005)
# (Cron version V5.0 -- $OpenBSD: crontab.c,v 1.48 2005/01/30 21:00:31 millert Exp $)
#minute hour    mday    month   wday    command
5	*	*	*	*	/usr/bin/find /home/will/Maildir -type d | /usr/bin/xargs -I\% chmod 700 \%
5	*	*	*	*	/usr/bin/find /home/will/Maildir -type f | /usr/bin/xargs -I\% chmod 600 \%
*/10	*	*	*	*	/home/will/bin/getmail
15	*/4	*	*	*	/home/will/Backups/mail.sh	
45	*	*	*	*	/usr/bin/crontab -l > /home/will/CVS/dotfiles/crontab/`/bin/hostname`
