# crontab/merkur.atekomi.net
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp//crontab.VeVSEX5865 installed on Fri Dec 16 08:32:53 2005)
# (Cron version V5.0 -- $OpenBSD: crontab.c,v 1.48 2005/01/30 21:00:31 millert Exp $)
##################  BEGIN HEADERS
# Filename	: `crontab -e`
# Use		: specifies programs and commands to be run on a periodic
#		  basis
# Author	: Will Maier <willmaier@ml1.net>
# Started	: 2005.05.01
# Vim		: vim: set tw=0:
# Updated	: 2005.08.25 14:57:34 -0500
##################  END HEADERS

#minute hour    mday    month   wday    command

# retrieve mail
*/10	*	*	*	*	/home/will/bin/getmail

# create backups of mail

#*/20	*	*	*	*	/home/will/Backups/mail.sh -h
30	4	*	*	1	/home/will/Backups/mail.sh -w
30	5	1	*	*	/home/will/Backups/mail.sh -m
30	6	1	1	*	/home/will/Backups/mail.sh -y

# clean out unused maildir folders
25	*/3	*	*	*	/home/will/bin/maildirs

# ensure reasonable permissions on mail
5	*	*	*	*	/usr/bin/find /home/will/Maildir -type d | /usr/bin/xargs -I\% chmod 700 \%
5	*	*	*	*	/usr/bin/find /home/will/Maildir -type f | /usr/bin/xargs -I\% chmod 600 \%

# back this crontab up
45	*	*	*	*	/usr/bin/crontab -l > /home/will/CVS/dotfiles/crontab/`/bin/hostname`

# dump my todo lists every morning
1	6	*	*	1-6	TDL_DATABASE=/home/will/.tdldb-cae /usr/local/bin/tdl list | /usr/bin/mail -s "TODO: CAE" will
1	6	*	*	1-6	TDL_DATABASE=/home/will/.tdldb-hep /usr/local/bin/tdl list | /usr/bin/mail -s "TODO: HEP" will
1	6	*	*	1-6	TDL_DATABASE=/home/will/.tdldb /usr/local/bin/tdl list | /usr/bin/mail -s "TODO: Personal" will
